# com.doga.cat-litter-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.doga.cat-litter-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/path/to/cat-litter-monitor/live_monitor.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/path/to/cat-litter-monitor</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HEADLESS</key>
        <string>true</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/path/to/cat-litter-monitor/logs/monitor_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/path/to/cat-litter-monitor/logs/monitor_stderr.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
